# supervisord.conf
[supervisord]
loglevel = info
nodaemon = true
identifier = supervisor

[program:web]
command=python manage.py runserver

[program:djangoq]
command=python manage.py qcluster
